# packaging/pkgsrc/net/radsecproxy/PLIST
@comment $NetBSD: PLIST,v 1.0 2007/11/02 10:45:13 kolla Exp $
man/man1/radsecproxy.1
man/man5/radsecproxy.conf.5
sbin/radsecproxy
share/doc/radsecproxy/README
share/doc/radsecproxy/AUTHORS
share/doc/radsecproxy/COPYING
share/doc/radsecproxy/INSTALL
share/doc/radsecproxy/LICENSE
share/doc/radsecproxy/ChangeLog
share/examples/rc.d/radsecproxy
share/examples/radsecproxy/radsecproxy.conf
share/examples/radsecproxy/radsec-dynsrv.sh
@dirrm share/examples/radsecproxy
@dirrm share/doc/radsecproxy
